# systemd/embedding-service.service
[Unit]
Description=Vector Memory Embedding Service
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.config/vector-memory/env
ExecStart=%h/.openclaw/workspace/vector_memory_venv/bin/python3 %h/.openclaw/workspace/vector_memory/embedding_service.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
