# systemd/reindex.timer
[Unit]
Description=Vector Memory Reindex Timer
Requires=embedding-service.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
