# systemd/reindex.service
[Unit]
Description=Vector Memory Reindex Service
Requires=embedding-service.service
After=embedding-service.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/vector-memory/env
ExecStart=%h/.openclaw/workspace/vector_memory_venv/bin/python3 %h/.openclaw/workspace/vector_memory/memory_reindex.py
